# scripts/launchd/com.hestia.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hestia.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/andrewroman117/hestia/scripts/hestia-watchdog.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/andrewroman117/hestia</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>HESTIA_NTFY_TOPIC</key>
        <string>hestia-trading-49ceebba</string>
    </dict>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/andrewroman117/hestia/logs/watchdog-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/andrewroman117/hestia/logs/watchdog-launchd.error.log</string>
</dict>
</plist>
