# scripts/launchd/com.hestia.trading-bots.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hestia.trading-bots</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/andrewroman117/hestia/.venv/bin/python</string>
        <string>-m</string>
        <string>hestia.trading.bot_service</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/andrewroman117/hestia</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/andrewroman117/hestia/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/andrewroman117/hestia</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/andrewroman117/hestia/logs/trading-bots.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/andrewroman117/hestia/logs/trading-bots.error.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
